# deploy/tangra-client.service
[Unit]
Description=Tangra Client Daemon - IPAM & LCM Agent
Documentation=https://github.com/go-tangra/go-tangra-client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tangra-client daemon \
    --config-file /etc/tangra-client/config.yaml
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow read access to DMI tables (board/memory info) and system info
ReadOnlyPaths=/sys/class/dmi /sys/firmware/dmi /proc
# Allow read/write to config and state directory
ReadWritePaths=/etc/tangra-client
# Allow read/write to certificate live directory
ReadWritePaths=/etc/tangra-client/live
# Allow self-update of the binary (selfupdate writes a temp file in the same directory)
ReadWritePaths=/usr/local/bin
# Allow package manager metadata refresh (apt, dnf/yum)
ReadWritePaths=/var/lib/apt/lists /var/cache/apt
ReadWritePaths=/var/cache/dnf /var/cache/yum /var/lib/dnf /var/lib/rpm

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tangra-client

[Install]
WantedBy=multi-user.target
