# deploy/tangra-client.service
[Unit]
Description=Tangra Client Daemon - IPAM & LCM Agent
Documentation=https://github.com/go-tangra/go-tangra-client
After=network-online.target
Wants=network-online.target
# Start rate-limiting belongs in [Unit].
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/tangra-client daemon \
    --config-file /etc/tangra-client/config.yaml
Restart=always
RestartSec=10

# The agent is a root host-administration daemon: it runs scripts and
# go-tangra-actions workflows that manage the system (package upgrades, file and
# service changes, etc.) and self-updates its own binary. Filesystem sandboxing
# (ProtectSystem/ProtectHome/ReadWritePaths/NoNewPrivileges/…) is intentionally
# NOT enabled — it makes the host unmanageable: e.g. `apt upgrade` fails with
# "Read-only file system" on /var/lib/dpkg, /var/cache, /usr, /boot, …

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tangra-client

[Install]
WantedBy=multi-user.target
